# crontab
0 5 1 * * /usr/bin/certbot renew --no-self-upgrade
* * * * * /dyndns.sh >/dev/null 2>&1
